# root/etc/services.d/codimd/run
#!/usr/bin/with-contenv bash

exec node $APP_PATH/app.js
